# infrastructure/ansible/collections/ansible_collections/cub_oit_pe/consul/roles/configure/templates/consul.service.j2
[Unit]
ConditionFileNotEmpty={{ configuration_location }}/consul.hcl

[Service]
ExecStart=/usr/bin/consul agent -config-dir={{ configuration_location }}